# data/dist/etc/veles/services.d/mn.discoveryDaemon.conf
[program:system.mn.discoveryDaemon]
command = /usr/libexec/veles/run-job DiscoveryDaemon
autorestart = true
stopasgroup = true
autostart = true
startsecs = 5
directory = /var/lib/veles
user = {{svc-user:system.net.webServer}}
