# data/dist/etc/veles/services.d/mn.metricDaemon.conf
[program:system.mn.metricDaemon]
command = /usr/libexec/veles/run-job MetricDaemon
autorestart = true
stopasgroup = true
autostart = true
startsecs = 5
directory = /var/lib/veles
user = {{svc-user:system.net.webServer}}
