# data/dist/etc/veles/services.d/net.tlsTunnel.conf
[program:system.net.tlsTunnel]
command = /etc/veles/services.d/net.tlsTunnel.start
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
directory = /
user = velesmn-stunnel
