# data/dist/etc/veles/services.d/vpn.obfsProxy.conf
[program:service.vpn.obfsProxy]
command = /usr/bin/obfsproxy --data-dir=/tmp/veles.scramblesuit scramblesuit --dest=127.0.0.1:{{svc-port:service.vpn.openVPNServer}} server 0.0.0.0:{{svc-port:service.vpn.obfsProxy}} --password SSSTIZ3LG5SSS43OMSSSM4LKGQFASSSS 
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
user = {{svc-user:service.vpn.obfsProxy}}
