# data/dist/etc/veles/services.d/vpn.dnsServer.conf
[program:service.vpn.dnsServer]
command = /usr/sbin/dnsmasq --conf-file=/etc/veles/vpn/dnsmasq.conf --no-daemon
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
#user = dns
