# data/dist/etc/veles/services.d/wallet.velesCoreDaemon.conf
[program:system.wallet.velesCoreDaemon]
command = /usr/bin/velesd -conf=/etc/veles/veles.conf -datadir=/var/lib/veles/wallet
autorestart = true
stopasgroup = true
autostart = true
startsecs = 30
user = {{svc-user:system.wallet.velesCoreDaemon}}
