# data/package/package.list
/etc/veles
/etc/systemd/system/veles-mn.service
/usr/libexec/veles
/usr/share/veles
/usr/bin/velesctl
/var/lib/veles/web
