# dist/etc/veles/services.d/db.redisServer.conf
[program:system.db.redisServer]
command = /usr/bin/redis-server /etc/veles/redisdb.conf --daemonize no
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
user = {{svc-user:system.db.redisServer}}
