# dist/etc/veles/services.d/net.tlsTunnel.conf
[program:system.net.tlsTunnel]
command = /usr/bin/stunnel4 /etc/veles/vpn/stunnel.conf
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
directory = /
user = {{svc-user:system.net.tlsTunnel}}
