# dist/etc/veles/services.d/net.protocolMux.conf
[program:system.net.protocolMux]
command = /usr/sbin/sslh -p 0.0.0.0:443 --openvpn 127.0.0.1:{{svc-port:service.vpn.openVPNServer}} --ssl 127.0.0.1:{{svc-port:system.net.tlsTunnel}}  --anyprot 127.0.0.1:{{svc-port:system.wallet.velesCoreDaemon}} -f
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
#user = sslh	# we need port 443
