# dist/etc/veles/services.d/net.protocolMux.conf
[program:system.net.protocolMux]
command = /usr/sbin/sslh -p 0.0.0.0:443 --openvpn 127.0.0.1:21341 --ssl 127.0.0.1:21339  --anyprot 127.0.0.1:21337 -f
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
#user = sslh	# we need port 443
