# dist/etc/veles/services.d/net.webServer.conf
[program:system.net.webServer]
command = /usr/bin/python3 /usr/lib/veles/web/webserver.py /etc/veles/webserver.conf
autorestart = true
stopasgroup = true
autostart = true
startsecs = 5
directory = /var/lib/veles/web
user = {{svc-user:system.net.webServer}}
