# dist/etc/veles/services.d/net.obfsProxy.conf
[program:system.net.obfsProxy]
command = /usr/bin/obfsproxy --data-dir=/tmp/veles.scramblesuit scramblesuit --dest=127.0.0.1:21341 server 0.0.0.0:21343 --password SSSTIZ3LG5SSS43OMSSSM4LKGQFASSSS 
autorestart = true
stopasgroup = true
autostart = true
startsecs = 10
#user = stunnel4	# we need to bing port 443
