# dist/etc/veles/services.d/vpn.openVPNServer.conf
[program:service.vpn.openVPNServer]
command = /usr/sbin/openvpn --config /etc/veles/vpn/openvpn.conf
autorestart = true
stopasgroup = true
autostart = true
startsecs = 15

